# tensin-app-gozim/gozim.sv.conf
[program:gozim]
command=/bin/bash -c "ZIM=$(ls -1 /zim/*.zim | tail -1 | cut -d'.' -f1) && ./gozimhttpd -path=${ZIM}.zim -index=${ZIM}.idx"
redirect_stderr=true
stdout_logfile=auto
