# tensin-app-gitblit/gitblit.sv.conf
[program:gitblit]
command=/bin/bash -c "java -server -Xmx1024M -Djava.awt.headless=true -jar /opt/gitblit/gitblit.jar --baseFolder /opt/gitblit-data"
redirect_stderr=true
stdout_logfile=auto
